# root/etc/profile.d/sma_path.sh
# env check
if [[ -z "${SMA_PATH}" ]]; then
    export SMA_PATH="/usr/local/sma"
fi
